# SqlInjection/docker/supervisord.conf
[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:mysql]
command = /usr/bin/mysqld_safe

[program:spring-boot]
command=java  -Djava.security.egd=file:/dev/./urandom -jar app.jar